# tmp/droid-config.files
%defattr(-,root,root,-)
/etc/acdbdata
/etc/audio_platform_info.xml
/etc/gps.conf
/etc/izat.conf
/etc/lowi.conf
/etc/media_codecs.xml
/etc/media_codecs_ffmpeg.xml
/etc/media_codecs_google_audio.xml
/etc/media_codecs_google_telephony.xml
/etc/media_codecs_google_video.xml
/etc/media_codecs_performance.xml
/etc/media_profiles.xml
/etc/mixer_paths.xml
/etc/modprobe.d/blacklist-evbug.conf
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/noplugins/ftp
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/mas
/etc/obexd/noplugins/syncevolution
/etc/obexd/opp.conf
/etc/obexd/pbap.conf
/etc/obexd/root
/etc/ofono/binder.conf
/etc/ofono/binder.d/dual-sim.conf
/etc/sap.conf
/etc/sec_config
/etc/selinux/config
/etc/sensorfw/primaryuse.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/Sensors.conf
/etc/xtra_root_cert.pem
/etc/xtwifi.conf
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/rules.d/998-droid-system.rules
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/droid-load-firmware.sh
/usr/bin/droid/droid-reclaim-memory.sh
/usr/bin/droid/get_partition.sh
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/sysctl.d/ipv6-settings.conf
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/droid-reclaim-memory.service
/usr/lib/systemd/system/droid-reclaim-memory.timer
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/timers.target.wants/droid-reclaim-memory.timer
/usr/lib/tmpfiles.d/droid-hal.conf
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/nemo/60-multimedia.conf
/var/lib/environment/ofono/noplugin.conf
/etc/audio_effects.conf
/etc/audio_effects_vendor.conf
/etc/audio_policy.conf
/etc/firmware/bcm43430A1.conf
/etc/gst-droid/gstdroidcamsrc-0.conf
/etc/gst-droid/gstdroidcamsrc-1.conf
/etc/gst-droid/gstdroidcodec.conf
/etc/media_codecs_google_video_le.xml
/etc/ofono/ril_subscription.conf
/etc/sensitive_pn.xml
/etc/sensorfw/primaryuse.conf
/etc/sensorfw/sensord.conf.d/60-sensors-grandneove3g.conf
/etc/sysconfig/bluetooth-rfkill-event-hciattach
/etc/tiny_hw.xml
/etc/audio_para
/etc/audio_hw.xml
/usr/bin/droid/btaddr.sh
/usr/bin/logcat
/usr/lib/systemd/system/network.target.wants/btaddr.service
/usr/lib/systemd/system/ofono.service.d/1-fixpts.conf
/usr/lib/systemd/system/btaddr.service
/usr/libexec/droid-hybris/system/etc/init/at_distributor.rc
/usr/libexec/droid-hybris/system/etc/init/atrace.rc
/usr/libexec/droid-hybris/system/etc/init/bootstat.rc
/usr/libexec/droid-hybris/system/etc/init/chown_service.rc
/usr/libexec/droid-hybris/system/etc/init/debuggerd.rc
/usr/libexec/droid-hybris/system/etc/init/drmserver.rc
/usr/libexec/droid-hybris/system/etc/init/dumpstate.rc
/usr/libexec/droid-hybris/system/etc/init/engpc.rc
/usr/libexec/droid-hybris/system/etc/init/gatekeeperd.rc
/usr/libexec/droid-hybris/system/etc/init/gpsd.rc
/usr/libexec/droid-hybris/system/etc/init/init-debug.rc
/usr/libexec/droid-hybris/system/etc/init/installd.rc
/usr/libexec/droid-hybris/system/etc/init/kill_phone.rc
/usr/libexec/droid-hybris/system/etc/init/logcatd.rc
/usr/libexec/droid-hybris/system/etc/init/logd.rc
/usr/libexec/droid-hybris/system/etc/init/macloader.rc
/usr/libexec/droid-hybris/system/etc/init/mdnsd.rc
/usr/libexec/droid-hybris/system/etc/init/mediacodec.rc
/usr/libexec/droid-hybris/system/etc/init/mediadrmserver.rc
/usr/libexec/droid-hybris/system/etc/init/mediaextractor.rc
/usr/libexec/droid-hybris/system/etc/init/modem_control.rc
/usr/libexec/droid-hybris/system/etc/init/modemd.rc
/usr/libexec/droid-hybris/system/etc/init/mtpd.rc
/usr/libexec/droid-hybris/system/etc/init/nvitemd.rc
/usr/libexec/droid-hybris/system/etc/init/perfprofd.rc
/usr/libexec/droid-hybris/system/etc/init/phoneserver.rc
/usr/libexec/droid-hybris/system/etc/init/racoon.rc
/usr/libexec/droid-hybris/system/etc/init/rild.rc
/usr/libexec/droid-hybris/system/etc/init/set_mac.rc
/usr/libexec/droid-hybris/system/etc/init/smd_symlink.rc
/usr/libexec/droid-hybris/system/etc/init/superuser.rc
/usr/libexec/droid-hybris/system/etc/init/uncrypt.rc
/usr/libexec/droid-hybris/system/etc/init/vdc.rc
/usr/libexec/droid-hybris/system/etc/init/vold.rc
/usr/libexec/droid-hybris/system/etc/init/data.rc
/usr/sbin/pts.sh
/var/lib/environment/compositor/droid-hal-device.conf
/usr/share/ssu/features.d/adaptation-community.ini
