# sparse/usr/lib/systemd/system/ofono.service.d/1-fixpts.conf
[Service]                                                                                                                                                             
ExecStartPre=/bin/sh /usr/sbin/pts.sh
